# package/DEBIAN/postinst
#!/bin/sh
set -e
systemctl enable pi-clock
systemctl start pi-clock
