# package/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop pi-clock || true
systemctl disable pi-clock || true
